# s6-overlay/s6-rc.d/otbr-web/finish
#!/command/with-contenv bash
echo "otbr-web ended with exit code ${1} (signal ${2})."
